# systemd/assistant.service
[Unit]
Description=Local Pi Assistant
After=network.target ollama.service
Wants=ollama.service

[Service]
Type=simple
User=simon
WorkingDirectory=/home/simon/Local-Pi-Assistant
ExecStart=/bin/bash /home/simon/Local-Pi-Assistant/start.sh
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
